# docker-files/docker/app/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=httpd -DFOREGROUND
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0